# build/debian/miru.socket
[Unit]
Description=Miru Config Agent socket

[Socket]
ListenStream=/run/miru/miru.sock
RuntimeDirectory=miru
SocketMode=0660
SocketGroup=miru
RemoveOnStop=yes

[Install]
WantedBy=sockets.target