# build/debian/postrm
#!/bin/sh

remove() {
    # reload the systemd daemon
    printf "\033[32m Reloading the systemd daemon\033[0m\n"
    systemctl daemon-reload

    # stop the systemd services
    printf "\033[32m Stopping the miru socket\033[0m\n"
    deb-systemd-invoke stop miru.socket

    # disable the systemd services
    printf "\033[32m Disabling the miru service\033[0m\n"
    deb-systemd-helper disable miru.service
    printf "\033[32m Disabling the miru socket\033[0m\n"
    deb-systemd-helper disable miru.socket
}

purge() {
    remove

    # reload the systemd daemon
    printf "\033[32m Reloading the systemd daemon\033[0m\n"
    systemctl daemon-reload

    # disable the systemd services
    printf "\033[32m Purging the miru service\033[0m\n"
    deb-systemd-helper purge miru.service
    printf "\033[32m Purging the miru socket\033[0m\n"
    deb-systemd-helper purge miru.socket

    # remove the files created by the miru agent
    printf "\033[32m Removing /var/lib/miru \033[0m\n"
    rm -rf /var/lib/miru
    printf "\033[32m Removing /var/log/miru \033[0m\n"
    rm -rf /var/log/miru
    printf "\033[32m Removing /srv/miru \033[0m\n"
    rm -rf /srv/miru
}

upgrade() {
    printf "\033[32m Post Remove of an upgrade\033[0m\n"
}

echo "$@"

action="$1"

case "$action" in
  "0" | "remove")
    printf "\033[32m Post Remove of a normal remove\033[0m\n"
    remove
    ;;
  "1" | "upgrade")
    upgrade
    ;;
  "purge")
    printf "\033[32m Post Remove purge\033[0m\n"
    purge
    ;;
  *)
    printf "\033[32m Alpine\033[0m"
    remove
    ;;
esac

exit 0