# build/debian/postinst
#!/bin/sh

# determine the version of systemd
socket_name="miru.socket"

create_miru_user() {
  if ! id -u miru > /dev/null 2>&1; then
    printf "\033[32m Creating the 'miru' user\033[0m\n"
    useradd -r -s /bin/false miru
  else
    printf "\033[32m The 'miru' user already exists\033[0m\n"
  fi
}

create_miru_group() {
  if ! getent group miru > /dev/null 2>&1; then
    printf "\033[32m Creating the 'miru' group\033[0m\n"
    groupadd -r miru
  fi
}

clean_install() {
    socket_name="$1"
    printf "\033[32m Post Install of an clean install\033[0m\n"

    # create the miru user and group
    create_miru_user
    create_miru_group

    # enable the socket
    printf "\033[32m Enable the socket\033[0m\n"
    systemctl enable "${socket_name}"

    # reload the unit from disk 
    printf "\033[32m Reload the service unit from disk\033[0m\n"
    systemctl daemon-reload

    # socket installation
    printf "\033[32m Set the enabled flag for the socket unit\033[0m\n"
    systemctl enable "${socket_name}"
    printf "\033[32m Restart the socket\033[0m\n"
    systemctl restart "${socket_name}"
}

upgrade() {
    socket_name="$1"
    printf "\033[32m Post Install of an upgrade\033[0m\n"

    # create the miru user and group
    create_miru_user
    create_miru_group

    # reload the unit from disk 
    printf "\033[32m Reload the service unit from disk\033[0m\n"
    systemctl daemon-reload

    # enable the socket
    printf "\033[32m Enable the socket\033[0m\n"
    systemctl enable "${socket_name}"

    # socket upgrade
    printf "\033[32m Restart the socket\033[0m\n"
    systemctl restart "${socket_name}"
}

# Step 2, check if this is a clean install or an upgrade
action="$1"
if [ "$1" = "configure" ] && [ -z "$2" ]; then
  action="install"
elif [ "$1" = "configure" ] && [ -n "$2" ]; then
  action="upgrade"
fi

case "$action" in
  "install")
    clean_install ${socket_name}
    ;;
  "upgrade")
    printf "\033[32m Post Install of an upgrade\033[0m\n"
    upgrade ${socket_name}
    ;;
esac

exit 0