# build/debian/postrm
#!/bin/sh

reload_systemd() {
    printf "\033[32m Reloading the systemd daemon\033[0m\n"
    systemctl daemon-reload
}

remove() {
    # stop the systemd services
    printf "\033[32m Stopping the miru socket\033[0m\n"
    deb-systemd-invoke stop miru.socket
    printf "\033[32m Stopping the miru service\033[0m\n"
    deb-systemd-invoke stop miru.service

    # disable the systemd services
    printf "\033[32m Disabling the miru service\033[0m\n"
    deb-systemd-helper disable miru.service
    printf "\033[32m Disabling the miru socket\033[0m\n"
    deb-systemd-helper disable miru.socket
}

purge() {
    remove

    # purge the systemd services
    printf "\033[32m Purging the miru service\033[0m\n"
    deb-systemd-helper purge miru.service
    printf "\033[32m Purging the miru socket\033[0m\n"
    deb-systemd-helper purge miru.socket

    # remove the files created by the miru agent
    printf "\033[32m Removing /var/lib/miru \033[0m\n"
    rm -rf /var/lib/miru
    printf "\033[32m Removing /var/log/miru \033[0m\n"
    rm -rf /var/log/miru
    printf "\033[32m Removing /srv/miru \033[0m\n"
    rm -rf /srv/miru
}

echo "$@"

action="$1"

reload_systemd

case "$action" in
  "0" | "remove")
    printf "\033[32m Post Remove of a normal remove\033[0m\n"
    remove
    ;;
  "1" | "purge")
    printf "\033[32m Post Remove purge\033[0m\n"
    purge
    ;;
  *)
    printf "\033[32m Post Remove of an unknown action\033[0m\n"
    remove
    ;;
esac

reload_systemd

exit 0