# build/deb/DEBIAN/preinst
#!/bin/bash
set -e  # Exit on error
version="0.6.5"
deb_pkg_name="miru"

miru_binaries_dir="/usr/local/miru"

agent_bin_name="miru-$version"
agent_bin_path="$miru_binaries_dir/$agent_bin_name"
agent_bin_symlink_name="miru"
agent_bin_symlink_path="$miru_binaries_dir/$agent_bin_symlink_name"
agent_start_script_name="miru-start.sh"

install_bin_name="install-miru"
install_bin_path="$miru_binaries_dir/$install_bin_name"

systemd_dir="/etc/systemd/system"
systemd_service_name="miru.service"
systemd_service_path="$systemd_dir/$systemd_service_name"

miru_dir="/var/lib/miru"
miru_log_dir="/var/log/miru"

# Variables
script="preinst"

# Ensure the script is running with elevated privileges
if [ "$(id -u)" -ne 0 ]; then
    echo "'$script' script must be run as the root user"
    exit 1
fi

# We can't uninstall the debian package since this script is run while installing the 
# package (debian can't do both at once)
if dpkg-query -W -f='${Status}' $deb_pkg_name 2>/dev/null | grep -q "install ok installed"; then
    echo ''
    echo "ERROR"
    echo "-----"
    echo "Miru is already installed. To reinstall please first uninstall $deb_pkg_name by running: sudo dpkg --purge $deb_pkg_name"
    echo ''
    echo ''
    exit 1
fi

# Stop, disable and remove the service if it exists
if systemctl list-unit-files | grep -q $systemd_service_name; then
    if systemctl is-active --quiet $systemd_service_name; then
        systemctl stop $systemd_service_name
    fi

    if systemctl is-enabled --quiet $systemd_service_name; then
        systemctl disable $systemd_service_name
    fi
fi

# Remove the symlink to the miru agent binary
if [ -L "$agent_bin_symlink_path" ]; then
    rm -f "$agent_bin_symlink_path"
fi

# Remove any binaries from a previous installation
if [ -d "$miru_binaries_dir" ]; then
    rm -rf "$miru_binaries_dir"
fi

# Remove any previous installation data
if [ -d "$miru_dir" ]; then
    rm -rf "$miru_dir"
fi

# Remove any previous log data
if [ -d "$miru_log_dir" ]; then
    rm -rf "$miru_log_dir"
fi

# Remove the service file if it exists
if [ -f "$systemd_service_path" ]; then
    rm -f "$systemd_service_path"
fi

exit 0