# build/deb/DEBIAN/prerm
#!/bin/bash
set -e # Exit on error
version="0.6.5"
deb_pkg_name="miru"

miru_binaries_dir="/usr/local/miru"

agent_bin_name="miru-$version"
agent_bin_path="$miru_binaries_dir/$agent_bin_name"
agent_bin_symlink_name="miru"
agent_bin_symlink_path="$miru_binaries_dir/$agent_bin_symlink_name"
agent_start_script_name="miru-start.sh"

install_bin_name="install-miru"
install_bin_path="$miru_binaries_dir/$install_bin_name"

systemd_dir="/etc/systemd/system"
systemd_service_name="miru.service"
systemd_service_path="$systemd_dir/$systemd_service_name"

miru_dir="/var/lib/miru"
miru_log_dir="/var/log/miru"

# Variables
script="prerm"

# Ensure the script is running with elevated privileges
if [ "$(id -u)" -ne 0 ]; then
  echo "'$script' script must be run as the root user"
  exit 1
fi

# Check if the service is running and stop it
if systemctl is-active $systemd_service_name 2>/dev/null --quiet; then
    systemctl stop $systemd_service_name 2>/dev/null --quiet
fi

# Disable the service to prevent it from starting on boot
if systemctl is-enabled $systemd_service_name 2>/dev/null --quiet; then
    systemctl disable $systemd_service_name 2>/dev/null --quiet
fi

exit 0