# build/deb/etc/systemd/system/miru.service
[Unit]
Description=Miru OTA deployment manager
After=network.target
Requires=miru.socket

[Service]
Type=simple
Group=miru
User=miru
ExecStart=/usr/local/bin/miru/config-agent

WorkingDirectory=/var/lib/miru
StandardOutput=journal
StandardError=journal
SyslogIdentifier=miru

[Install]
WantedBy=multi-user.target
