# build/deb/etc/systemd/system/miru.socket
[Unit]
Description=Miru OTA deployment manager socket

[Socket]
ListenStream=/run/miru.sock
SocketMode=0660
SocketGroup=miru
RemoveOnStop=yes

[Install]
WantedBy=sockets.target